# pi-server/config/coachmate-sync.service
[Unit]
Description=CoachMate schedule sync (Supabase -> local cache)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/pi/route-tracker/pi-server
ExecStart=/usr/bin/node sync-schedule.mjs
User=pi
# SUPABASE_ENV=prod is the default (see sync-schedule.mjs) — set to dev here
# only while testing against the dev project.
# Environment=SUPABASE_ENV=prod

[Install]
WantedBy=multi-user.target
